# mirror/hdctools/beaglebone/servod.conf
# Copyright 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start servod"
author          "chromium-os-dev@chromium.org"

start on net-device-up
stop on runlevel [!2345]
respawn
# If you change the respawn limit, change the log rotation, below.
respawn limit 3 10

script
  # Only port 9999 is used on beaglebone devices.
  LOG="/var/log/servod_9999.STARTUP.log"
  {
    . /etc/default/servod
    if [ -n "$BOARD" ]; then
      BOARD_FLAG="--board ${BOARD}"
      BOARD_MSG="Launching servod for ${BOARD}"
      logger -t "$UPSTART_JOB" "$BOARD_MSG"
      echo "$BOARD_MSG"
    else
      BOARD_FLAG=""
      NO_BOARD_MSG="Launching servod with unspecified DUT board"
      logger -t "$UPSTART_JOB"
      echo "$NO_BOARD_MSG"
    fi

    echo "================= Version information ================="
    echo "$REVISION"
    echo "======================================================="
  } > $LOG 2>&1
  exec servod --log-dir --debug --host 0.0.0.0 $BOARD_FLAG
end script
