# script/from_rd/v2/lcfc_stress_0104/warm-boot-lcfc.conf
# Copyright (c) 2013 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "reboot"
author          "chromium-os-dev@chromium.org"

start on starting system-services

script
  counter=`cat /usr/local/lcfc_stress_0104/wb_count`
  echo $(($counter+1)) > /usr/local/lcfc_stress_0104/wb_count
  sleep 1
  if [ $counter -lt 2500 ]; then
    date >> /var/log/reboot.log
    sleep 90
    echo $(($counter+1)) >> /var/log/reboot.log
    echo "reboot" >> /var/log/reboot.log
    powerd_setuid_helper --action=reboot
    #reboot
  fi
end script
